# services/bytebin/bytebin-service.conf
[program:bytebin]
command=/usr/bin/java -jar bytebin.jar
directory=/opt/bytebin
autostart=true
autorestart=true
stderr_logfile=/var/log/bytebin.err.log
stdout_logfile=/var/log/bytebin.out.log